# etc/systemd/system/suppress-gpe4E.service
[Unit]
Description=Disable GPE 4E

[Service]
ExecStart=/usr/bin/bash -c 'echo "disable" > /sys/firmware/acpi/interrupts/gpe4E'

[Install]
WantedBy=multi-user.target 
